# app/systemd/bsky-firehose.service
[Unit]
Description=Bluesky Firehose Listener
After=network-online.target

[Service]
Type=simple
User=bskybots
Group=bskybots
EnvironmentFile=/etc/bsky-bots.env
WorkingDirectory=/opt/bsky-bots
ExecStart=/opt/bsky-bots/.venv/bin/python /opt/bsky-bots/bsky-firehose.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/bsky-bots/firehose.log
StandardError=append:/var/log/bsky-bots/firehose.err

[Install]
WantedBy=multi-user.target
